# docker/php-apache/crontab
* * * * * cd /var/www/html/ && ./cron.sh
